# etc/profile.d/texlive.sh
export PATH=${PATH:+${PATH}:}/usr/local/texlive/2019/bin/x86_64-linux
export INFOPATH=${INFOPATH:+${INFOPATH}:}/usr/local/texlive/2019/texmf-dist/doc/info
