# etc/profile.d/texlive.sh
export PATH=/usr/local/texlive/texlive-dist/bin/x86_64-linux${PATH:+:$PATH}
export INFOPATH=/usr/local/texlive/texlive-dist/texmf-dist/doc/info${INFOPATH:+:$INFOPATH}
